# docker/daiquiri/rootfs/home/dq/conf/crontab
# min hour dom mon dow command
0 3 * * * ${HOME}/sh/archive-jobs.sh
0 3 * * * ${HOME}/sh/cleanup-download-files.sh
